# root/etc/services.d/celery-beat/run
#!/usr/bin/with-contenv sh
cd /app/api
exec s6-setuidgid funkwhale celery -A funkwhale_api.taskapp beat --pidfile= -f /var/log/funkwhale/celery-beat.log
