# root/etc/services.d/celery-beat/run
#!/usr/bin/with-contenv sh
cd $FUNKWHALE_PATH/api
exec s6-setuidgid funkwhale celery -A funkwhale_api.taskapp beat
