# root/etc/services.d/redis/run
#!/usr/bin/with-contenv sh
redis-server
