# root/etc/services.d/postgres/run
#!/usr/bin/with-contenv
exec s6-setuidgid postgres postgres -D /var/lib/postgresql/data
