# root/etc/services.d/postgres/run
#!/usr/bin/with-contenv sh
exec s6-setuidgid funkwhale postgres -D /data/data >> /var/log/postgresql/postgresql.log 2>&1
