# root/etc/services.d/celery-worker/run
#!/usr/bin/with-contenv sh
cd /app/api
exec s6-setuidgid funkwhale celery -A funkwhale_api.taskapp worker --concurrency=$CELERYD_CONCURRENCY -f /var/log/funkwhale/celery-worker.log
