# root/etc/services.d/funkwhale/run
#!/usr/bin/with-contenv sh
set -a && . /defaults/pre_env && . "$FUNKWHALE_PATH/config/.env" && set +a
cd /srv/funkwhale/api
exec s6-setuidgid funkwhale daphne config.asgi:application
