# root/etc/services.d/funkwhale/run
#!/usr/bin/with-contenv sh
cd /app/api
exec s6-setuidgid funkwhale daphne config.asgi:application >> /var/log/funkwhale/funkwhale.log 2>&1
