# root/etc/services.d/funkwhale/run
#!/usr/bin/with-contenv sh
cd /srv/funkwhale/api
export $(awk '!/^(#|$)/ { print $1 }' .env)
s6-setuidgid funkwhale daphne config.asgi:application
