# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
mkdir -p /run/nginx
rm /etc/nginx/conf.d/default.conf
nginx -g "daemon off;"
